# admindb/crontab
0 3 * * * root /bin/backup >> /var/log/cron.log 2>&1

